# buildroot/board/OpenLF/Didj/file_delete.list
etc/init.d/S05avahi-setup.sh
etc/init.d/S40network
