# board/OpenLF/Didj/files/etc/profile.d/sdl_exports.sh
#!/bin/sh
# SDL exports
export SDL_NOMOUSE=1
export SDL_FBDEV="/dev/fb0"
export SDL_VIDEODRIVER="fbcon"
